# debian/media-hub.user.service
[Unit]
Description=Starts the media-hub service

[Service]
ExecStart=/usr/bin/media-hub-server

[Install]
WantedBy=default.target
